# templates/default/systemd/squid.service.erb
# Dynamically generated by Chef - local modifications will be overwritten

[Unit]
Description=Squid Web Proxy Cache Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Environment="SQUID_ARGS=-YC"
ExecStartPre=<%= node['squid']['binary'] %> -z -f <%= node['squid']['config_file'] %>
ExecStart=<%= node['squid']['binary'] %> -N $SQUID_ARGS -f <%= node['squid']['config_file'] %>
ExecStop=<%= node['squid']['binary'] %> -k shutdown
ExecReload=<%= node['squid']['binary'] %> -k reconfigure
KillSignal=SIGINT
TimeoutStopSec=40

[Install]
WantedBy=multi-user.target
